# tahoe-night-nurse.service
[Unit]
Description=Tahoe Night Nurse Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/app
EnvironmentFile=-/root/app/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tahoe-night-nurse

[Install]
WantedBy=multi-user.target
